# infrastructure/ansible/roles/oxidize/templates/oxidize.service.j2
[Unit]
Description=Oxidize Relay Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/oxidize
ExecStart=/opt/oxidize/oxidize-server --config /etc/oxidize/server.toml
Restart=always
RestartSec=3
LimitNOFILE=2097152
LimitMEMLOCK=infinity
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_IPC_LOCK
Environment=RUST_LOG=info
Environment=RUST_BACKTRACE=1
EnvironmentFile=-/etc/oxidize/auth.env

[Install]
WantedBy=multi-user.target
