# deploy/systemd/oxidize.service
[Unit]
Description=Oxidize Network Relay Server
After=network.target
Documentation=https://github.com/yourusername/oxidize

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/oxidize
ExecStart=/home/ubuntu/oxidize/target/release/oxidize-server \
  --listen 0.0.0.0:4433 \
  --config /home/ubuntu/oxidize/production.toml \
  --metrics-addr 0.0.0.0:9090
  
Restart=always
RestartSec=10
StartLimitInterval=0

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oxidize

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
