# daemon/oxidize-daemon.service
[Unit]
Description=Oxidize Network Relay Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/oxidize-daemon
Restart=on-failure
RestartSec=1
Environment=RUST_LOG=info

# === ZERO-DOWNTIME DEPLOYMENT ===
# Use SIGTERM for graceful shutdown (30s drain timeout)
KillSignal=SIGTERM
# Allow time for connection draining before force kill
TimeoutStopSec=35

# Reload support - send SIGHUP to reload config without restart
ExecReload=/bin/kill -HUP $MAINPID

# Resource limits for high performance
LimitNOFILE=65535
LimitNPROC=4096

# Network capabilities for TUN
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN

# Security hardening (relaxed for packet capture)
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/var/run/oxidize

[Install]
WantedBy=multi-user.target
